# backend/services/ctl-001-workspace-identity/systemd/ndsp-ctl-001-workspace-identity.service
[Unit]
Description=NDSP CTL-001 Workspace Identity
After=network.target

[Service]
Type=simple
User=nawaf511
WorkingDirectory=/home/nawaf511/empire-core-new/backend/services/ctl-001-workspace-identity
Environment=NODE_ENV=production
Environment=NDSP_PORT=9081
Environment=NDSP_HOST=127.0.0.1
EnvironmentFile=-/home/nawaf511/empire-core-new/backend/.env
ExecStart=/usr/bin/node /home/nawaf511/empire-core-new/backend/services/ctl-001-workspace-identity/main.cjs
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
